# NexBridgeServer
#!/bin/sh
# /etc/init.d/NexBridgeServer
#
# NEXBRIDGE server startup script
# loops<> 
#
### BEGIN INIT INFO
# Provides:          nexbridge server
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: nexbridge server
### END INIT INFO

# DO NOT EDIT THIS FILE
# set your drivers in /usr/local/etc/NexBridgeServer.conf

DAEMON="/usr/local/bin/nexbridge"
LOCKFILE="/var/lock/nexbridge"
RUNAS="root"

SCRIPT="$0"
if [ -e '/usr/local/etc/NexBridgeServer.conf' ]; then
  . /usr/local/etc/NexBridgeServer.conf
else
  exit 1
fi

# set lock file
if [ ! -e $LOCKFILE ]; then
  touch $LOCKFILE
fi

test -f $DAEMON || exit 0

case "$1" in
  start)
	lock=`cat $LOCKFILE`
    	if [ !$lock ] || [ $lock = 0 ]; then
	echo "Starting NexBridge server"
      	su - $RUNAS -c "$DAEMON $DAEMON_OPTION > /dev/null 2>&1 &"
	pid=$(pidof nexbridge)
	echo $pid > $LOCKFILE
	else
	   echo "NexBridge server already running"
	fi
	;;
  stop)
  	lock=`cat $LOCKFILE`
	if [ $lock ] && [ $lock != 0 ]; then
	echo "Stopping NexBridge server"
      	kill -TERM $lock &>/dev/null
      	echo 0 > $LOCKFILE
    	else
      	   echo "No NexBrdige server running"
    	fi
    	;;

  restart)
    $SCRIPT stop && $SCRIPT start
    ;;
  status)
	lock=`cat $LOCKFILE`
	if [ $lock ] && [ $lock != 0 ] && [ `netstat -tlp|grep $lock|wc -l` = 1 
]; then
	echo "NexBridge server is running on $LISTEN:$PORT"
	else
	   echo "NexBridge server is NOT running"
	fi	
    	;;
  *)
    echo "Usage: /etc/init.d/NexBridgeServer {start|stop|status}"
    exit 1
    ;;
esac

exit 0
